# templates/docker-healthcheck.timer.j2
[Unit]
Description=Docker healthcheck timer for %I

[Timer]
OnCalendar={{ __time }}
Persistent=true
RandomizedDelaySec=5min
Unit=docker-healthcheck@%i.service

[Install]
WantedBy=timers.target